# scripts/deb/voip-opus-update.timer
# Poll for new voip-opus releases on a schedule.
#
#   OnBootSec=5min          first check 5 minutes after boot
#   OnUnitActiveSec=30min   then every 30 minutes
#   RandomizedDelaySec=5min jitter so a fleet of Pis doesn't all hit the
#                           GitHub API at the same wall-clock minute
#
# 30 minutes is a reasonable trade-off: snappy enough that users get fixes
# within an hour of release, infrequent enough that we stay well under the
# 60 req/hr unauthenticated GitHub API limit even with multiple Pis on the
# same network egress IP.

[Unit]
Description=Periodically check for voip-opus updates
Wants=voip-opus-update.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=5min
Persistent=true
Unit=voip-opus-update.service

[Install]
WantedBy=timers.target
